# Library/LaunchAgents/com.sweedler.ari.dotfiles.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sweedler.ari.dotfiles.plist</string>

  <key>Program</key>
  <string>/Users/ari.sweedler/.config/bin/git-dotfiles-fetch</string>

  <key>StandardErrorPath</key>
  <string>/Users/ari.sweedler/Desktop/git-dotfiles-fetch.err</string>

  <key>StandardOutPath</key>
  <string>/Users/ari.sweedler/Desktop/git-dotfiles-fetch.out</string>

  <key>StartInterval</key>
  <integer>86400</integer>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
