# Library/LaunchAgents/com.github.domt4.homebrew-autoupdate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.github.domt4.homebrew-autoupdate</string>

  <key>Program</key>
  <string>/Users/ari.sweedler/Library/Application Support/com.github.domt4.homebrew-autoupdate/brew_autoupdate</string>

  <key>ProgramArguments</key>
  <array>
      <string>/Users/ari.sweedler/Library/Application Support/com.github.domt4.homebrew-autoupdate/brew_autoupdate</string>
  </array>

  <key>StandardErrorPath</key>
  <string>/Users/ari.sweedler/Library/Logs/com.github.domt4.homebrew-autoupdate/com.github.domt4.homebrew-autoupdate.out</string>

  <key>StandardOutPath</key>
  <string>/Users/ari.sweedler/Library/Logs/com.github.domt4.homebrew-autoupdate/com.github.domt4.homebrew-autoupdate.out</string>

  <key>StartInterval</key>
  <integer>86400</integer>

  <key>LowPriorityBackgroundIO</key>
  <true/>

  <key>LowPriorityIO</key>
  <true/>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
