# Library/LaunchAgents/org.git-scm.git.daily.plist
<?xml version="1.0"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key>
  <string>org.git-scm.git.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/opt/git/libexec/git-core/git</string>
    <string>--exec-path=/opt/homebrew/opt/git/libexec/git-core</string>
    <string>for-each-repo</string>
    <string>--config=maintenance.repo</string>
    <string>maintenance</string>
    <string>run</string>
    <string>--schedule=daily</string>
  </array>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Day</key><integer>1</integer>
      <key>Hour</key><integer>0</integer>
      <key>Minute</key><integer>49</integer>
    </dict>
    <dict>
      <key>Day</key><integer>2</integer>
      <key>Hour</key><integer>0</integer>
      <key>Minute</key><integer>49</integer>
    </dict>
    <dict>
      <key>Day</key><integer>3</integer>
      <key>Hour</key><integer>0</integer>
      <key>Minute</key><integer>49</integer>
    </dict>
    <dict>
      <key>Day</key><integer>4</integer>
      <key>Hour</key><integer>0</integer>
      <key>Minute</key><integer>49</integer>
    </dict>
    <dict>
      <key>Day</key><integer>5</integer>
      <key>Hour</key><integer>0</integer>
      <key>Minute</key><integer>49</integer>
    </dict>
    <dict>
      <key>Day</key><integer>6</integer>
      <key>Hour</key><integer>0</integer>
      <key>Minute</key><integer>49</integer>
    </dict>
  </array>

  <key>StandardErrorPath</key>

  <key>StandardErrorPath</key>
  <string>/Users/ari.sweedler/Library/Logs/org.git-scm.git.daily/stderr.txt</string>

  <key>StandardOutPath</key>
  <string>/Users/ari.sweedler/Library/Logs/com.github.domt4.homebrew-autoupdate/com.github.domt4.homebrew-autoupdate.out</string>
</dict>
</plist>
